# build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
nocleanup=true

[program:dnsfilter]
command=/venv/bin/python3 -u /usr/src/app/src/dnsfilter/main.py
directory=/usr/src/app
user=named
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_events_enabled=false
stderr_events_enabled=false
stdout_capture_maxbytes=0
stderr_capture_maxbytes=0
redirect_stderr=true
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0